# deployment/playbooks/file_manager/generated/develop/file_manager.service
[Unit]
Description=file manager service

[Service]
Type=simple
WorkingDirectory=/root/deploy/file_manager
ExecStart=/root/deploy/file_manager/file_manager -config config/develop.json -port=443  1>/dev/null 2>&1 &
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
