# deployment/playbooks/file_manager/generated/local/file_manager.service
[Unit]
Description=file manager service

[Service]
Type=simple
WorkingDirectory=/home/vagrant/deploy/file_manager
ExecStart=/home/vagrant/deploy/file_manager/file_manager -config config/local.json -port=443  1>/dev/null 2>&1 &
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
