# deployment/playbooks/file_manager/generated/gcp-develop/file_manager.service
[Unit]
Description=file manager service

[Service]
Type=simple
WorkingDirectory=/home/jello-backend/deploy/file_manager
ExecStart=/home/jello-backend/deploy/file_manager/file_manager -config config/gcp-develop.json -port=7890  1>/dev/null 2>&1 &
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
